# deploy/systemd/hadlink-redirect.service
[Unit]
Description=hadlink URL Shortener Redirect Daemon
After=network.target
Documentation=https://github.com/jbsco/hadlink

[Service]
Type=simple
ExecStart=/usr/local/bin/hadlink redirect
User=hadlink
Group=hadlink
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadOnlyPaths=/var/lib/hadlink
MemoryMax=64M
CPUQuota=20%
Restart=on-failure
RestartSec=5s

# Environment
Environment=HADLINK_STORAGE=/var/lib/hadlink/hadlink.db
Environment=HADLINK_MODE=redirect
Environment=HADLINK_PORT=8080

[Install]
WantedBy=multi-user.target
